# systemd/opencode-web.service
[Unit]
Description=OpenCode Web Interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/opencode-web/env
Environment=BROWSER=true
ExecStart=%h/.local/share/mise/shims/opencode web --hostname 0.0.0.0 --port 4096
WorkingDirectory=%h
Restart=on-failure
RestartSec=5

# Memory limit
MemoryMax=8G
MemoryHigh=7G

# Security hardening
# NOTE: 不設 ProtectHome — opencode 是 coding 工具，user 需從 web UI 編輯
# 任意 project；而 service 跑在 user 自己 uid，ProtectHome 對同 uid 無意義。
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-web

[Install]
WantedBy=default.target
